# avaloq-acpr/src/ansible/roles/avaloq.acpr-docker/templates/dockercleanup.service.j2
[Unit]
Description=Docker housekeeping

[Service]
Type=oneshot
ExecStart={{ acpr_docker_cleanup_base }}/bin/dockercleanup.sh
